# django.service
# Contents of /etc/systemd/system/django.service
[Unit]
Description=Django
After=network.target

[Service]
Type=simple
Restart=on-failure
ExecStart=su -c "python /home/pi/rasberry-pi-camera-live/app/manage.py runserver 0.0.0.0:8080" pi

# Hardening
SystemCallArchitectures=native
MemoryDenyWriteExecute=true
NoNewPrivileges=true

[Install]
WantedBy=default.target